# obsidian/.config/systemd/user/obsidian-backup.timer
[Unit]
Description=Backup do Obsidian a cada 30 minutos

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
